# config/cron/production/crontab
######## Miscellaneous tasks ###########

*/30 * * * * /projects/scout/current/config/cron/production/rake.sh subscriptions:reinitialize

######## Delivery tasks ########

# email digests - once a day, in the morning
0 10 * * * /projects/scout/current/config/cron/production/deliver.sh email_daily

# immediate emails - throughout the day, every 10 minutes
*/10 * * * *  /projects/scout/current/config/cron/production/deliver.sh email_immediate

# SMS messages - throughout the day, every 10 minutes
*/10 * * * * /projects/scout/current/config/cron/production/deliver.sh sms


######## Polling tasks  #########

# keyword search subscriptions

# bills in Congress - daily, bounded by THOMAS/GPO/GovTrack
0 8,14 * * *  /projects/scout/current/config/cron/production/check.sh federal_bills

# bills in the states - once an hour, Open States scrapers running throughout the day, but no need to be real time
0 11,15 * * *     /projects/scout/current/config/cron/production/check.sh state_bills

# speeches in Congress - daily, bounded by GPO
0 8,14 * * *  /projects/scout/current/config/cron/production/check.sh speeches

# regulations - once an hour, FederalRegister.gov updates at unknown intervals
0 * * * *     /projects/scout/current/config/cron/production/check.sh regulations


# item-specific subscriptions

# activity on bills in Congress - daily, bounded by THOMAS/GovTrack
0 8,14 * * *  /projects/scout/current/config/cron/production/check.sh federal_bills_activity

# upcoming notices for bills in Congress - every 5 minutes, bounded by party leadership
*/5 * * * *   /projects/scout/current/config/cron/production/check.sh federal_bills_upcoming_floor

# upcoming hearings for bills in Congress - every 5 minutes, updated throughout the day
*/5 * * * *   /projects/scout/current/config/cron/production/check.sh federal_bills_hearings

# vote results for bills in Congress - every 5 minutes, updated by House/Senate throughout the day
*/5 * * * *   /projects/scout/current/config/cron/production/check.sh federal_bills_votes

# activity on bills in the States - once an hour, who knows when these will be updated
0 * * * * /projects/scout/current/config/cron/production/check.sh state_bills_activity

# activity on bills in the States - once an hour, who knows when these will be updated
0 * * * * /projects/scout/current/config/cron/production/check.sh state_bills_votes

# external feed subscriptions, once an hour
0 * * * * /projects/scout/current/config/cron/production/check.sh feed


# keep the newline below
