# config/cron/example/crontab
# SCOUT_PATH is the path to the clone of the repository
SCOUT_PATH=/path/to/repository
# SCOUT_ENV_FILE contains all environment variables for Scout
SCOUT_ENV_FILE=/path/to/.profile
# SCOUT_LOG_PATH keeps log files generated by Scout
SCOUT_LOG_PATH=/path/to/logs

# Assumes that server time is UTC and than local time is EST.

# Initialize any uninitialized subscriptions every half hour
*/30 *  * * * . $SCOUT_ENV_FILE && cd $SCOUT_PATH && rake subscriptions:reinitialize     > $SCOUT_LOG_PATH/subscriptions-reinitialize.log 2>&1

# At midnight, email a report of Google's crawling activity
0    5  * * * . $SCOUT_ENV_FILE && cd $SCOUT_PATH && rake analytics:google               > $SCOUT_LOG_PATH/analytics-google.log 2>&1

# At 12:15am on Monday, email about last week's user activity
15   5  * * 1 . $SCOUT_ENV_FILE && cd $SCOUT_PATH && rake analytics:weekly service=scout > $SCOUT_LOG_PATH/analytics-weekly.log 2>&1

# At 12:30am, regenerate the sitemap (you should sync before running this job)
30   5  * * * . $SCOUT_ENV_FILE && cd $SCOUT_PATH && rake sitemap                        > $SCOUT_LOG_PATH/sitemap.log 2>&1

### Delivery tasks

# At 8am, deliver daily emails
0    13 * * * . $SCOUT_ENV_FILE && cd $SCOUT_PATH && rake deliver:email_daily            > $SCOUT_LOG_PATH/deliver-email_daily.log 2>&1

# Every 10 minutes, deliver immediate emails
*/10 *  * * * . $SCOUT_ENV_FILE && cd $SCOUT_PATH && rake deliver:email_immediate        > $SCOUT_LOG_PATH/deliver-email_immediate.log 2>&1

### Polling tasks

# Every hour, check RSS and Atom feeds for new items
45   *  * * * . $SCOUT_ENV_FILE && cd $SCOUT_PATH && rake subscriptions:check:feed       > $SCOUT_LOG_PATH/check-feed.log 2>&1

# # Add any other syncing and polling jobs, depending on your subscription adapters
